# src/main/assembly/dist/install/dd-virus-scan.service
[Unit]
Description=Dd Virus Scan Service
After=payara.service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-virus-scan/bin/dd-virus-scan server /etc/opt/dans.knaw.nl/dd-virus-scan/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143
Environment="JAVA_OPTS=-Xms64m"
Environment="LC_ALL=en_US.UTF-8"

User=dd-virus-scan
Group=dd-virus-scan

[Install]
WantedBy=multi-user.target
